# systemd/dnews-push.timer.tmpl
[Unit]
Description=Daily News - Push Timer

[Timer]
{{PUSH_ONCALENDAR_LINES}}
AccuracySec=1s
Unit=dnews-push.service

[Install]
WantedBy=timers.target
